# docker/s6-overlay/s6-rc.d/horizon/run
#!/command/execlineb -P
with-contenv
php /var/www/html/artisan horizon
